# busybox-rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
ttySAC0::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
